# systemd/yt-abs-importer-worker.service
[Unit]
Description=yt-abs-importer Background Worker
After=network.target redis-server.service yt-abs-importer-app.service
Requires=redis-server.service

[Service]
Type=simple
User=yt-abs-importer
Group=yt-abs-importer
WorkingDirectory=/opt/yt-abs-importer
EnvironmentFile=/etc/yt-abs-importer/.env

# Fail gracefully with a clear error message if config is missing
ExecStartPre=/bin/sh -c 'if [ ! -f /etc/yt-abs-importer/.env ]; then echo "ERROR: /etc/yt-abs-importer/.env is missing. Please copy /opt/yt-abs-importer/.env.example and populate it."; exit 1; fi'

ExecStart=/opt/yt-abs-importer/.venv/bin/python -m rq.cli worker ytabs
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
